# crontab

# m h  dom mon dow   command

* * * * *  /root/update.sh
* * * * *  /root/update_dtrip/action.sh

